# cookbooks/growthforecast/templates/default/growthforecast.supervisord.ini.erb
[program:growthforecast]
command=<%= @perl_bin %> <%= @growthforecast_bin %> --data-dir <%= @data_dir %>
process_name=%(program_name)s
stdout_logfile_maxbytes=1MB
stderr_logfile_maxbytes=1MB
stdout_logfile=/var/log/supervisor/%(program_name)s.log
stderr_logfile=/var/log/supervisor/%(program_name)s.log
user=<%= @exec_user %>
autostart=true
autorestart=true
